# raspberryPi/config/raspberryPi.service
[Unit]
After=graphical.target

[Service]
ExecStart=/home/chaes/mechatronicSem5ProjCheese/raspberryPi/main.py

[Install]
WantedBy=default.target